# backend/systemd/fedora-website-generator.service
[Unit]
Description=Fedora AI Website Generator
Documentation=https://github.com/fedora-ai/website-generator
After=network.target redis.service
Wants=redis.service

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/opt/fedora-ai-website-generator
Environment=PATH=/opt/fedora-ai-website-generator/venv/bin
EnvironmentFile=/opt/fedora-ai-website-generator/.env
ExecStart=/opt/fedora-ai-website-generator/venv/bin/python backend/main.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5

# Security
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/fedora-ai-website-generator/uploads /opt/fedora-ai-website-generator/generated

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fedora-ai-website-generator

[Install]
WantedBy=multi-user.target